# services/zerotier-one/run
#!/bin/sh -e

cd $(readlink -f ${0} | xargs dirname)

ls /dev/net/tun > /dev/null

exec zerotier-one
