# services/zerotier-one/run
#!/bin/sh -e

cd $(readlink -f ${0} | xargs dirname)

TUN="/dev/net/tun"
if [ ! -f ${TUN} ]; then
    mkdir -p $(dirname ${TUN})
    mknod -m a+w ${TUN} c 10 200
fi

exec zerotier-one
